# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python3 apis.py
stderr_logfile=api_out.log
stdout_logfile=api_error.log
